# deploy/systemd/primal.service
[Unit]
Description=Primal SoundCloud listen-along
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=primal
Group=primal
WorkingDirectory=/opt/primal/current
EnvironmentFile=/etc/primal/primal.env
ExecStart=/usr/bin/node dist-server/server/index.js
Restart=on-failure
RestartSec=3s
TimeoutStopSec=15s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/primal
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=false
CapabilityBoundingSet=
SystemCallArchitectures=native
UMask=0027
LimitNOFILE=8192

[Install]
WantedBy=multi-user.target

